# launch_agents/com.gabriel.mem0-apple-intelligence.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gabriel.mem0-apple-intelligence</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>node</string>
        <string>/Volumes/Ready500/DEVELOPMENT/mem0/integrations/mcp/native_node_server.js</string>
    </array>
    
    <!-- Auto-start and keep alive -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Restart behavior -->
    <key>ThrottleInterval</key>
    <integer>10</integer>    
    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>QDRANT_URL</key>
        <string>http://localhost:6333</string>
        <key>QDRANT_COLLECTION</key>
        <string>gabriel_apple_intelligence_memories</string>
        <key>APPLE_INTELLIGENCE_ENABLED</key>
        <string>true</string>
        <key>LOG_LEVEL</key>
        <string>INFO</string>
    </dict>
    
    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/mem0-apple-intelligence.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mem0-apple-intelligence-error.log</string>
    
    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Volumes/Ready500/DEVELOPMENT/mem0</string>
    
    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>1024</integer>
    </dict>
</dict>
</plist>